# templates/node_exporter.service.epp
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/node_exporter <% if $t7d_prometheus::node_exporter::textfile == true { -%> --collector.textfile.directory="/var/lib/prometheus/node-exporter" <% } -%>

[Install]
WantedBy=multi-user.target
