# templates/varnish_exporter.service.epp
[Unit]
Description=Prometheus varnish Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=varnish_exporter
Group=varnish
Type=simple
ExecStart=/usr/local/bin/varnish_exporter
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
